# deployment/a434.ini
[program:a434]
command=/home/a434/venv/bin/gunicorn -b localhost:8181 -w 4 a434:app
directory=/home/a434
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true